# supervisor-mikrotik-worker.conf
[program:mikrotik-worker]
process_name=%(program_name)s_%(process_num)02d
command=php %(ENV_APP_BASE_PATH%/artisan queue:work database --queue=mikrotik --timeout=300 --sleep=3 --tries=3
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=%(ENV_USERNAME)s
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/worker.log
stopwaitsecs=3600
